# docker/wis2box-api.cron
*/5 * * * * python3 /app/docker/clean-old-jobs.py > /proc/1/fd/1 2>/proc/1/fd/2
